# systemd/vpn-portal.service
[Unit]
Description=VPN portal login gate (form auth) + monitor
After=network.target

[Service]
User=www-data
Group=www-data
ExecStart=/usr/bin/python3 /opt/vpn-portal/app.py
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadOnlyPaths=/opt/vpn-portal /etc/vpn-portal

[Install]
WantedBy=multi-user.target
